# advanced-monitoring/Worker/crontab
59 23 * * * rm -f /##WORKINGDIR##/*.csv
